# scripts/com.catalyst.local.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Catalyst Local — macOS LaunchAgent.

  Runs the production `next start` build on port 3030 so the dashboard is
  reachable from localhost even when Vercel is suspended. Auto-starts on
  login, restarts on crash (throttled), bounded RAM and FDs.

  Install (copies this plist into ~/Library/LaunchAgents and loads it):
      scripts/local-daemon.sh install

  This file is the source of truth; the daemon script copies it to
  ~/Library/LaunchAgents/com.catalyst.local.plist on install.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.catalyst.local</string>

  <!-- Direct invocation of pnpm — no shell wrapper so we don't hit TCC
       restrictions when launchd tries to chdir into ~/Desktop. PATH is
       set explicitly via EnvironmentVariables below so child processes
       (node, etc.) resolve. -->
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/pnpm</string>
    <string>--dir</string>
    <string>/Users/manuelgarrido/dev/catalyst-local</string>
    <string>start</string>
    <string>--port</string>
    <string>3030</string>
  </array>

  <!-- WorkingDirectory omitted intentionally — letting launchd start the
       process in /tmp and using `pnpm --dir` to target the repo avoids
       the TCC "Operation not permitted" error that fires when a
       LaunchAgent tries to chdir into a Desktop-protected directory.
       The user only needs to grant /opt/homebrew/bin/pnpm permission
       (or accept the per-folder prompt the first time it reads
       ~/Desktop/catalyst-local). -->

  <key>RunAtLoad</key>
  <true/>

  <!-- KeepAlive with SuccessfulExit=false means: relaunch the process if
       it exits with a non-zero status. Plain crashes get respawned, but
       a clean shutdown doesn't loop. -->
  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>

  <!-- Wait 30 seconds between respawn attempts so a broken build doesn't
       spin and burn CPU. With launchd's default minimum (10s) plus this,
       a crash loop is bounded. -->
  <key>ThrottleInterval</key>
  <integer>30</integer>

  <!-- Resource limits applied to the child process. Subidos 2026-05-23
       tras EAGAIN -11 en ESM reads: launchd default soft limit (256) era
       muy bajo para Next 16 + React 19 production que abre ~600 fds
       durante el bootstrap. SoftResourceLimits 8192 / Hard 16384 deja
       margen 10× sin riesgo (sistema típicamente cap a 200k). -->
  <key>SoftResourceLimits</key>
  <dict>
    <key>NumberOfFiles</key>
    <integer>8192</integer>
  </dict>
  <key>HardResourceLimits</key>
  <dict>
    <key>NumberOfFiles</key>
    <integer>16384</integer>
  </dict>

  <!-- LOCAL_MODE + LOCAL_DEFAULT_SESSION_ID let the server fall back to a
       known session UUID when the localhost cookie is missing, so the
       user's watchlist appears without manual cookie injection. Empty
       string here = no fallback; set in the user's shell or .env.local
       to enable. -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>NODE_ENV</key>
    <string>production</string>
    <key>LOCAL_MODE</key>
    <string>1</string>
    <!-- --max-old-space-size removido 2026-05-23: con ese flag (cualquier
         valor) el ESM loader sync-read crasheaba con EAGAIN -11 al cargar
         chunks de la app desde launchd (la causa real resultó ser iCloud
         File Provider en ~/Desktop; el repo vive ahora en ~/dev). NO
         reintroducir el flag de heap: V8 default ~1.7GB y Next 16 ronda
         ~250-400MB en idle.
         max-http-header-size SÍ es necesario: Yahoo Finance responde
         headers >16KB (CSP+cookies) y undici los rechaza con
         UND_ERR_HEADERS_OVERFLOW, rompiendo /api/article on-click
         (2026-07-18). Flag no relacionado con heap-sizing. -->
    <key>NODE_OPTIONS</key>
    <string>--max-http-header-size=65536</string>
    <!-- Explicit PATH so pnpm finds node + spawned binaries even though
         we don't go through a login shell.
         IMPORTANT: /opt/homebrew/opt/node@22/bin must precede the default
         /opt/homebrew/bin because Node 23.6.1 (current Homebrew default)
         has an ESM loader regression on macOS that crashes `next start`
         with EAGAIN -11 in readFileHandle. node@22 LTS does not. When
         brew default catches up to 23.7+ or 24 LTS the prefix can drop.
         See user memory: feedback_node_23_macos_eagain. -->
    <key>PATH</key>
    <string>/opt/homebrew/opt/node@22/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOME</key>
    <string>/Users/manuelgarrido</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/manuelgarrido/dev/catalyst-local/.next/daemon-logs/stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/manuelgarrido/dev/catalyst-local/.next/daemon-logs/stderr.log</string>

  <!-- Throttle CPU priority so the daemon never starves foreground apps
       on the user's Mac. -->
  <key>ProcessType</key>
  <string>Background</string>
  <key>Nice</key>
  <integer>5</integer>
</dict>
</plist>
